# profile.d/206-path-nvim.sh
# add the nvim inary to the path
PATH="$PATH:/opt/nvim-linux64/bin"
